# ns_wan_gain.service
[Unit]
Description=Nazaré-coupled WAN egress gain daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /home/martin/claude/ns_wan_gain.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
