# pi2/beacon-pi2.service
[Unit]
Description=Kuramoto beacon sid=2 (pi2)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/beacon pi2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
